# srcs/grafana/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:grafana]
command=sh -c "./grafana-server && kill $(pidof supervisord)"
autostart=true
autorestart=true